# init/librarian.conf
description "Resource Librarian Slackbot"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
chdir /var/docker-services/librarian
exec docker-compose up
